# unzip_bot.service
[Unit]
Description=Telegram Unzip Bot
After=network.target

[Service]
Type=simple
User=YOUR_USERNAME
WorkingDirectory=/path/to/Unzip_File_Bot
ExecStart=/usr/bin/python3 bot.py
Restart=always
RestartSec=10
StandardOutput=append:/var/log/unzip_bot.log
StandardError=append:/var/log/unzip_bot_error.log

[Install]
WantedBy=multi-user.target
